# workers/launchd/fm.hoku.constellation.detect-changes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>fm.hoku.constellation.detect-changes</string>
    <key>ProgramArguments</key>
    <array>
        <string>__REPO_DIR__/workers/run.sh</string>
        <string>detect-changes</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__REPO_DIR__/workers/logs/detect-changes.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/workers/logs/detect-changes.err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
